# etc/profile.d/scrssi.sh
alias scrssi="screen -raAd -S scrssi && clear"
